# infrastructure/ansible/roles/manager/templates/swarmd-manager.service.j2
[Unit]
Description=SwarmCracker SwarmKit Manager
Documentation=https://github.com/restuhaqza/SwarmCracker
After=network.target network-online.target docker.service consul.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart={{ swarmcracker_install_dir }}/swarmd-firecracker \
  --manager \
  --state-dir {{ manager_state_dir }} \
  --listen-control-api /var/run/swarmkit/swarm.sock \
  --hostname {{ ansible_hostname }} \
  --listen-remote-api {{ manager_bind_address }}:{{ manager_port }} \
  --advertise-remote-api {{ manager_advertise_address }}:{{ manager_advertise_port }} \
  --kernel-path {{ swarmcracker_kernel_path }} \
  --rootfs-dir {{ swarmcracker_rootfs_dir }} \
  --bridge-name {{ swarmcracker_bridge_name }} \
  --subnet {{ swarmcracker_subnet }} \
  --bridge-ip {{ swarmcracker_bridge_ip }} \
  --nat-enabled \
  --enable-cni \
  --cni-plugin-dir {{ cni_plugin_dir | default('/opt/cni/bin') }} \
  --cni-config-dir {{ cni_config_dir | default('/etc/cni/net.d') }} \
  --cni-subnet-pool {{ cni_subnet_pool | default('10.0.0.0/8') }} \
  --cni-subnet-size {{ cni_subnet_size | default(24) }} \
  --cni-vxlan-port {{ cni_vxlan_port | default(4789) }} \
  {% if swarmcracker_vxlan_enabled %}--vxlan-enabled{% endif %} \
  {% if consul_enabled | default(true) %}--consul-enabled --consul-address {{ consul_client_addr | default('127.0.0.1') }}:8500{% endif %} \
  {% if swarmcracker_debug %}--debug{% endif %}

ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -SIGTERM $MAINPID

Restart=on-failure
RestartSec=10s
TimeoutStopSec=30s

# Environment
Environment=SWARM_SOCKET=/var/run/swarmkit/swarm.sock
Environment=SWARM_STATE_DIR={{ manager_state_dir }}

# Security hardening
# ReadWritePaths removed to avoid mount namespace issues

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target